# gemini-live.service
[Unit]
Description=Gemini Live Voice Server
After=network.target caddy.service
Wants=caddy.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/experiments/gemini-live
EnvironmentFile=/root/.openclaw/workspace/experiments/gemini-live/.env
ExecStart=/root/.openclaw/workspace/experiments/gemini-live/.venv/bin/python server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gemini-live

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/.openclaw/workspace/experiments/gemini-live /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
